# setup/autostart/XfceFvwm.desktop
[Desktop Entry]
Exec=fvwm --replace

# or
# Exec=fvwm3 --replace
