# setup/autostart/XfceFvwm2-3.desktop
[Desktop Entry]
Exec=fvwm --replace
Exec=fvwm3 --replace
